# files/flask/flaskgunicorn.service.j2
[Unit]
Description=flask gunicorn daemon
After=network.target

[Service]
User=chrisotelios
Group=www-data
WorkingDirectory={{service.workingdir}}
ExecStart={{service.execstart}}

[Install]
WantedBy=multi-user.target
